# api.conf
[program:api]
command=/usr/bin/python2.7 apimain.py
directory=/opt/raspberrypid/backend
user=root
autostart=true
autorestart=true