# backend.conf
[program:backend]
command=/usr/bin/nice -n 15 /usr/bin/python2.7 main.py
directory=/opt/raspberrypid/backend
user=root
autostart=true
autorestart=true